# profile.d/20oos_oraclexe.sh
export ORACLE_SID=XE
export ORAENV_ASK=NO
. ORACLE_HOME/bin/oraenv
alias sqlplus='rlwrap sqlplus'
alias rman='rlwrap rman'
